# deploy/abay-alerts.service
[Unit]
Description=ABAY Alert Monitor
After=network.target redis-server.service

[Service]
Type=simple
User=abay
Group=abay
WorkingDirectory=/home/abay/abay-app/django_backend
Environment="PATH=/home/abay/abay-app/venv/bin"
EnvironmentFile=/home/abay/abay-app/.env
ExecStart=/home/abay/abay-app/venv/bin/python \
    manage.py monitor_alerts --interval 60
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
